# infra/systemd/user/vh-analysis-backend-3001.service
[Unit]
Description=VHC Analysis Backend (:3001 analyze health/config + article-text)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Override VHC_REPO and ExecStart in a drop-in if your checkout lives elsewhere.
Environment=VHC_REPO=/home/humble/VHC
Environment=PATH=/home/humble/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOST=127.0.0.1
Environment=PORT=3001
Environment=ARTICLE_TEXT_MAX_CHARS=24000
Environment=ARTICLE_FETCH_TIMEOUT_MS=12000
ExecStart=/usr/bin/env node /home/humble/VHC/tools/scripts/vh-analysis-backend-3001.js
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
