# infra/systemd/user/vh-news-aggregator.service
[Unit]
Description=VHC News Aggregator Publisher Daemon
Documentation=file:/home/humble/VHC/docs/ops/news-aggregator-production-service.md
After=network-online.target vh-storycluster-engine.service
Wants=network-online.target vh-storycluster-engine.service

[Service]
Type=simple
# Override VHC_REPO, WorkingDirectory, ExecStart, and VH_NEWS_DAEMON_ENV_FILE in a drop-in if the checkout or env file differs.
Environment=VHC_REPO=/home/humble/VHC
Environment=VH_NEWS_DAEMON_ENV_FILE=%h/.config/vhc/news-aggregator.env
Environment=PATH=/home/humble/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/humble/VHC
ExecStart=/usr/bin/env bash /home/humble/VHC/tools/scripts/start-news-aggregator-daemon-production.sh
Restart=always
RestartSec=30
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
