# infra/systemd/user/vh-phase5-scope-a-watch-closure.timer
[Unit]
Description=Build VHC Phase 5 Scope A watch closure packet every 30 minutes

[Timer]
OnBootSec=12min
OnUnitActiveSec=30min
AccuracySec=2min
Persistent=true
Unit=vh-phase5-scope-a-watch-closure.service

[Install]
WantedBy=timers.target
