# tools/launchd/com.vhc.automation-stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vhc.automation-stack</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd __REPO_ROOT__ &amp;&amp; exec bash tools/scripts/automation-stack.sh ensure</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME_DIR__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/.tmp/automation-stack/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/.tmp/automation-stack/logs/launchd-stderr.log</string>
</dict>
</plist>
